# infra/transcript-vault-mcp/scripts/sync.timer.example
[Unit]
Description=Run transcript-vault nightly sync at 03:00 local time

[Timer]
OnCalendar=*-*-* 03:00:00

# If the laptop was asleep / powered off when 03:00 hit, run as soon
# as we wake up. This is the whole reason we're using systemd over cron.
Persistent=true

# Smear the start across a 5-minute window in case multiple timers
# converge on 03:00.
RandomizedDelaySec=5min

# Keep AccuracySec tight so we actually fire near the requested time.
AccuracySec=1min

Unit=transcript-vault-sync.service

[Install]
WantedBy=timers.target
